# perp-backend.service
[Unit]
Description=Perp Backend Docker Compose Stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/zoopx@corp.zoopx.xyz/perp-backend
RemainAfterExit=yes
Environment=COMPOSE_PROJECT_NAME=perp-backend
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
